# crontab
0 0 * * * . /app/cron.env && cd /app && python src/run_pipeline.py >> /var/log/pipeline.log 2>&1
